# audio_client.service
[Unit]
Description=Audio client
After=multi-user.target

[Service]
User=pi
Type=simple

WorkingDirectory=/home/pi/radio_control
ExecStart=/home/pi/radio_control/audio_client.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
